# rover_service/rover.service
[Unit]
After=network.service

[Service]
User=ubuntu
Group=ubuntu
ExecStart=/bin/bash /usr/local/bin/rover.sh

[Install]
WantedBy=default.target
