# rover_service/rover.service
[Unit]
After=network.service

[Service]
ExecStart=/usr/local/bin/rover.sh

[Install]
WantedBy=default.target
